# templates/etc/cron.daily/mapcrafter.j2
#!/bin/bash

/usr/local/bin/mapcrafter -c /etc/mapcrafter.conf
